# trunk/osinfo/init.d/osinfo.gentoo
#!/sbin/runscript
# osinfo init script for Gentoo Linux
# Copyright 2006 Mikael Lammentausta
# Distributed under the terms of the BSD License

source /etc/osinfo.conf

depend() {
	after *
}

start() {
	ebegin "Running osinfo"
	start-stop-daemon --start --exec /usr/local/bin/osinfo -- "${osinfo_options}" "${osinfo_modules}" -o "${osinfo_output_file}"
	eend $? 
}

